# airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
Type=simple
Environment=XDG_SESSION_TYPE=wayland
ExecStart=
ExecStart=-/sbin/agetty --autologin liveuser --noclear %I 38400 linux
